# helperpoc-backup.app/Contents/Library/LaunchDaemons/com.keypath.helperpoc.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Launch daemon plist for SMAppService privileged helper -->
<!-- BundleProgram uses app-bundle relative path (required for SMAppService) -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" 
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.keypath.helperpoc.helper</string>
    <key>BundleProgram</key>
    <string>helperpoc-helper</string>
    <key>AssociatedBundleIdentifiers</key>
    <array>
        <string>com.keypath.helperpoc.helper</string>
    </array>
    <key>MachServices</key>
    <dict>
        <key>com.keypath.helperpoc.xpc</key>
        <true/>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>